# scripts/wispr-thoughts.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Wispr Thoughts Sunday 8am local weekly digest. Bootstrap fills the {{...}} placeholders. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{LABEL}}</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "{{INSTALL_PATH}}" && /Library/Frameworks/Python.framework/Versions/3.12/bin/python3 src/weekly_email.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{INSTALL_PATH}}</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/wispr-thoughts.out.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/wispr-thoughts.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{HOME}}/.npm-global/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>WISPRTHOUGHTS_LLM_BACKEND</key>
        <string>api</string>
    </dict>
</dict>
</plist>
